# 0-HGES/Scripts/.config/autostart/ConfigureVNC.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Type=Application
Name=Configure-VNC
Icon=
Exec=/usr/share/HGeS/configuraVNC.sh
