# 0-HGES/Scripts/.config/autostart/lembrete.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Type=Application
Name=lembrete
Icon=
Exec=/usr/share/HGeS/lembrete.sh
