# 0-HGES/Scripts/.config/autostart/user-first-config.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Type=Application
Name=user-first-config
Icon=
Exec=/usr/share/HGeS/user-first-config.sh
