# 3-9BLOG/2 - LINUXMINT_21.2/home/administrador/.config/autostart/lembrete.desktop
#!/usr/bin/env xdg-open

[Desktop Entry]
Type=Application
Name=lembrete
Icon=1E64_notepad.0
Exec=/usr/share/9BLOG/lembrete.sh
